# container-files/etc/supervisor.d/nginx.conf
[program:nginx]
command = /usr/local/sbin/nginx -g 'daemon off;'
autorestart = true
stderr_logfile = /var/log/nginx/error.log
stdout_logfile = /var/log/nginx/access.log

# Watch for changes in Nginx conf directories and restart Nginx when a config change occured
[program:nginx-reload]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /etc/nginx/; do nginx -t && nginx -s reload; done'
